# etc/crontab
# m h  dom mon dow   command
0 5  * * *  /home/tigaron/docker/server.fls.red/cronjob.sh
